# supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
command=/bin/bash -c "service postgresql start"
